# kuavo-ros-control/tools/udev_rules/90-roban-audio.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0d8c", ATTR{idProduct}=="013c", RUN+="/sbin/modprobe snd-usb-audio"